# snapshots/systemd/subs-check.service.d/limits.conf
# Resource limits for subs-check
# 防止节点测速吃光 RAM 影响其他服务 (OpenClaw / Hermes / Web UI)
[Service]
MemoryMax=512M
MemoryHigh=384M
CPUQuota=80%
